# root/etc/s6-overlay/s6-rc.d/init-mods-universal-acl-add-package/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
printf "acl\n" >> /mod-repo-packages-to-install.list
